# dri/systemd/penalty-partner-rsi.service
[Unit]
Description=Penalty Partner autonomous DRI (Tibo) — one RSI cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Shadow-first: SHIP=0 = observe+decide+report, NO auto merge/deploy until graduated.
Environment=PP_DRI_SHIP=0
Environment=TIBO_POST_HOURS=01 09 17
ExecStart=/bin/bash /root/penalty-partner/dri/rsi.sh
TimeoutStartSec=2600
