# dri/systemd/penalty-partner-watchdog.service
[Unit]
Description=Penalty Partner DRI heartbeat watchdog

[Service]
Type=oneshot
ExecStart=/bin/bash /root/penalty-partner/dri/watchdog.sh
TimeoutStartSec=120
