# dri/systemd/penalty-partner-rsi.timer
[Unit]
Description=Hourly Penalty Partner DRI (Tibo) RSI loop

[Timer]
OnCalendar=*-*-* *:17:00
Persistent=true

[Install]
WantedBy=timers.target
