# dri/systemd/push-or-pay-rsi.service
[Unit]
Description=Push or Pay autonomous DRI (Tibo) — one RSI cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PP_DRI_SHIP=0
Environment=TIBO_POST_HOURS=01 09 17
ExecStart=/bin/bash /root/push-or-pay/dri/rsi.sh
TimeoutStartSec=2600
