# dri/systemd/push-or-pay-watchdog.timer
[Unit]
Description=Push or Pay DRI watchdog (offset from RSI loop)
[Timer]
OnCalendar=*-*-* *:47:00
Persistent=true
[Install]
WantedBy=timers.target
